# build/package/mita/debian/DEBIAN/postinst
#!/bin/sh

/usr/bin/id mita > /dev/null 2>&1
rc=$?
if [ $rc -ne 0 ]; then
    /usr/sbin/useradd --no-create-home --user-group mita
fi

set -e

mkdir -p /etc/mita
chmod 775 /etc/mita

systemctl daemon-reload

# Server daemon will run with the system.
systemctl enable mita.service
systemctl start mita.service
